# openrc/mbimapidle.in
#!/sbin/openrc-run

description="Execute command on IMAP mailbox changes"
command="@BINDIR@/mbimapidle"
supervisor=supervise-daemon

